# services/memory-sync-service.service
[Unit]
Description=Memory Sync Service - REST API for remote memory synchronization
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=lj
Group=lj
WorkingDirectory=/home/lj/haivemind-mcp-server
Environment=PYTHONPATH=/home/lj/haivemind-mcp-server/src
Environment=HAIVEMIND_JWT_SECRET=cfccc737f768f7f8de4d5436da6535da5ca3f5f26df715978f4bd2f51932c68e
ExecStart=/home/lj/haivemind-mcp-server/venv/bin/python /home/lj/haivemind-mcp-server/src/sync_service.py --host 0.0.0.0 --port 8899
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/lj/haivemind-mcp-server/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target