# services/memory-mcp-server.service
[Unit]
Description=Memory MCP Server - Distributed memory storage for Claude
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=lj
Group=lj
WorkingDirectory=/home/lj/haivemind-mcp-server
Environment=PYTHONPATH=/home/lj/haivemind-mcp-server/src
ExecStart=/home/lj/haivemind-mcp-server/venv/bin/python /home/lj/haivemind-mcp-server/src/memory_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/lj/haivemind-mcp-server
PrivateTmp=true

[Install]
WantedBy=multi-user.target